# deploy/kore-batch-sample.service
[Unit]
Description=Kore Batch Sample
After=network.target postgresql.service

[Service]
Type=simple
User=kore-batch
Group=kore-batch

WorkingDirectory=/opt/kore-batch

# Secrets dans un fichier separe - jamais dans ce fichier
# Creer /etc/kore-batch/env depuis deploy/env.prod.template
EnvironmentFile=/etc/kore-batch/env

ExecStart=/usr/bin/java \
    -Xmx512m \
    -jar /opt/kore-batch/kore-batch-sample.jar \
    --inputFile=${INPUT_FILE}

StandardOutput=journal
StandardError=journal
SyslogIdentifier=kore-batch-sample

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
